# pirage.conf
[program:pirage]
command=/home/bjp/.pyenv/versions/3.6.4/bin/python /home/bjp/pirage/run.py
directory=/home/bjp/pirage
autostart=true
autorestart=true
