# pirage.conf
[program:pirage]
command=/home/pi/venvs/pirage/bin/python /home/pi/pirage/run.py
directory=/home/pi/pirage
autostart=true
autorestart=true
