# pirage/pirage.conf
[program:pirage]
command=/home/pi/venv/
directory=/home/pi/pirage
autostart=true
autorestart=true
